# udev/70-steamos-cec-toolkit.rules
SUBSYSTEM=="cec", TAG+="uaccess", RUN+="/var/lib/steamos-cec-toolkit/steamos-cec-permissions-apply"
