# systemd/user/cec-audio-control.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/.local/bin/steamos-cec-external-volume
